# udev/99-drm-switchmode.rules
ACTION=="change" SUBSYSTEM=="drm" RUN+="/usr/local/sbin/video_switch.sh"


